# ansible/roles/vm_monitoring_agent/templates/otel-collector.service.j2
[Unit]
Description=Aether OpenTelemetry Collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ monitoring_agent_service_user }}
ExecStart={{ otel_collector_binary_path }} --config={{ otel_collector_config_path }}
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aether-otel-collector
# Grant access to system logs and metrics
PrivateDevices=no
ProtectSystem=no
ProtectHome=no
ReadWritePaths={{ otel_collector_storage_path }}
ReadOnlyPaths=/var/log

[Install]
WantedBy=multi-user.target
