# ansible/playbooks/openbao/templates/openbao-cert-renew.service.j2
[Unit]
Description=Renew OpenBao TLS certificate from step-ca
After=network-online.target openbao.service
Wants=network-online.target
BindsTo=openbao.service

[Service]
Type=simple
User=root
Group=root

Environment="STEPPATH={{ openbao_config_dir }}/step"

# Renewal script re-splits certs for IAM Roles Anywhere and restarts OpenBao
ExecStart=/usr/bin/step ca renew --daemon \
  --exec "{{ openbao_config_dir }}/renew-hook.sh" \
  {{ openbao_tls_dir }}/openbao.crt \
  {{ openbao_tls_dir }}/openbao.key

Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target

