# ansible/playbooks/openbao/templates/openbao-cert-renew.service.j2
[Unit]
Description=Renew OpenBao TLS certificate from step-ca
After=network-online.target openbao.service
Wants=network-online.target
Requires=openbao.service

[Service]
Type=simple
User={{ openbao_user }}
Group={{ openbao_group }}

Environment="STEPPATH={{ openbao_config_dir }}/step"

ExecStart=/usr/bin/step ca renew --daemon \
  --exec "systemctl reload openbao.service || systemctl restart openbao.service" \
  {{ openbao_tls_dir }}/openbao.crt \
  {{ openbao_tls_dir }}/openbao.key

Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target

