# ansible/roles/technitium_dns/templates/technitium-dns-health-gate.service.j2
# Managed by Ansible (technitium_dns role). Do not edit by hand.
[Unit]
Description=Technitium DNS health gate
After=network-online.target technitium-dns.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/technitium-dns-health-gate
