# ansible/roles/technitium_dns/templates/technitium-dns-health-gate.timer.j2
# Managed by Ansible (technitium_dns role). Do not edit by hand.
[Unit]
Description=Run Technitium DNS health gate every ten seconds

[Timer]
OnBootSec=10s
OnUnitActiveSec=10s
AccuracySec=1s
Persistent=true
Unit=technitium-dns-health-gate.service

[Install]
WantedBy=timers.target
